# provision/claude-telegram-tmux.service
# Always-on Claude Code + Telegram running inside an ATTACHABLE tmux session.
# Installed by `kapp-service-install` (the default mode). Use this when you want
# to WATCH Claude work live (`tmux attach`) while still getting auto-restart and
# start-on-boot. The __DEVUSER__ placeholder is replaced with the real non-root
# user at install time.
#
#   Watch live:  tmux attach -t claudebot      (detach: Ctrl+B then D)
#   Logs:        journalctl -u claude-telegram
#
# How it stays alive: a `while` loop inside tmux relaunches Claude ~15s after it
# exits/crashes (claude-telegram-run.sh waits for setup, then execs claude). So
# crash/OOM recovery comes from the loop; systemd just starts the tmux session on
# boot and tears it down on stop. (For pure systemd supervision with no tmux, use
# `kapp-service-install --headless`, which installs claude-telegram.service.)

[Unit]
Description=KAppMaker AI — always-on Claude + Telegram (attachable tmux)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=__DEVUSER__
# HOME + a fixed tmux socket dir so `tmux attach` from your login shell finds it.
Environment=HOME=/home/__DEVUSER__
Environment=TMUX_TMPDIR=/tmp
# systemd does NOT read ~/.bashrc — load the toolchain env so bot-driven Gradle/
# Android builds find JAVA_HOME / ANDROID_SDK_ROOT / PATH. '-' = ok if missing.
EnvironmentFile=-/home/__DEVUSER__/.config/kappmaker/env
# If the box runs low on memory, prefer killing a heavy build JVM over this bot.
OOMScoreAdjust=-500
ExecStart=/usr/bin/tmux new-session -d -s claudebot "while true; do /home/__DEVUSER__/bin/claude-telegram-run.sh; echo '[claude-telegram] session ended — restarting in 15s'; sleep 15; done"
ExecStop=-/usr/bin/tmux kill-session -t claudebot

[Install]
WantedBy=multi-user.target
