# provision/setup-web.service
# KAppMaker AI — one-time browser-setup terminal (ttyd) for new boxes.
# Installed by provision/bootstrap.sh; __DEVUSER__ is replaced at install time.
#
# Serves provision/setup-wizard.sh in a web terminal so a non-technical customer
# can finish setup (Claude login + Telegram token) from the browser — no SSH.
#
# Security model:
#   - ttyd binds LOOPBACK ONLY (-i lo); the outside world reaches it through
#     Caddy, which terminates real TLS on an sslip.io hostname (see bootstrap.sh).
#   - Access is gated by HTTP basic auth: user "setup", password = the random
#     SETUP_CODE generated at bootstrap (shown to the customer in the dashboard).
#     Anyone with the code gets a shell as the box owner — acceptable because
#     the code is only ever shown to the box's owner and the whole box IS theirs.
#   - EnvironmentFile has no "-" prefix ON PURPOSE: if the code file is missing
#     (e.g. after teardown) the unit must fail rather than start with no auth.
#   - Once setup completes, setup-teardown.sh stops+disables this unit, closes
#     the firewall ports, and writes the marker below so it can never come back.

[Unit]
Description=KAppMaker AI — one-time browser setup terminal (ttyd)
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/etc/kappmaker/.setup-done

[Service]
Type=simple
User=__DEVUSER__
Environment=HOME=/home/__DEVUSER__
WorkingDirectory=/home/__DEVUSER__
EnvironmentFile=/etc/kappmaker/setup-web.env
# -i lo: loopback only (Caddy fronts it)   -W: allow typing (ttyd 1.7+ is
# read-only by default)   -m 2: at most two concurrent viewers.
ExecStart=/usr/bin/ttyd --port 7681 --interface lo --writable --max-clients 2 \
  --credential setup:${SETUP_CODE} \
  -t titleFixed=KAppMaker-Setup -t disableLeaveAlert=true \
  /home/__DEVUSER__/bin/setup-wizard.sh
Restart=always
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
